# systemd/user/redshift.service
[Unit]
Description=Redshift display colour temperature adjustment

[Service]
ExecStart=/usr/bin/redshift
Restart=always

[Install]
WantedBy=default.target
